# claude-code-development-kit/scripts/com.magpie.devkit-cron.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.magpie.devkit-cron</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__PLUGIN_ROOT__/scripts/devkit-cron.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PLUGIN_ROOT__</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>__STATE_DIR__/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__STATE_DIR__/launchd-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>NO_COLOR</key>
        <string>1</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
